# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery Notification Service

[Service]
ExecStart=/home/cl1ngz/Documents/Scripts/battery-notify.sh
