# .config/systemd/user/rclone-sync.timer
[Unit]
Description=Run rclone sync once a month

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
